# dist/twentyfouronline-scheduler.service
[Unit]
Description=twentyfouronline Scheduler

[Service]
Type=oneshot
StandardOutput=null
StandardError=null
WorkingDirectory=/opt/twentyfouronline/
ExecStart=/usr/bin/env php artisan schedule:run --no-ansi --no-interaction
User=twentyfouronline
Group=twentyfouronline




